# src/sd-taskpacker/sd-taskpacker.service
[Unit]
Description=sd-taskpacker

[Service]
BusName=org.verbum.SdTaskPacker
ExecStart=/usr/bin/sd-taskpacker
